# deployment/debian/postinst
#!/bin/bash

if id gnomeshade &>/dev/null; then
	echo 'user gnomeshade exists'
else
	echo 'user gnomeshade does not exist, creating'
	adduser --system --no-create-home --home /nonexistent --disabled-login --shell=/bin/false gnomeshade || exit
fi

chown gnomeshade /etc/opt/gnomeshade --recursive || exit
chown gnomeshade /opt/gnomeshade --recursive || exit

systemctl daemon-reload || exit
systemctl enable gnomeshade.service || exit
