# deployment/debian/prerm
#!/bin/bash

systemctl stop gnomeshade.service || exit
systemctl disable gnomeshade.service || exit
