# deployment/debian/gnomeshade/DEBIAN/postrm
#!/bin/bash
set -e

if [[ $(systemctl) =~ -\.mount ]]
then
	exit
else
	update-rc.d gnomeshade remove
fi
